# for-crontab/poll-pull-git-refresh-dock-container.cron
*/2 * * * * /usr/bin/sh /base/dock-script/reinit-vps-gateway-access.sh
*/2 * * * * /usr/bin/sh /base/web-script/git-pull-restart-container.sh
