# for-crontab/poll-pull-git-refresh-dock-container.cron
*/2 * * * * ubuntu sh /home/ubuntu/va-sta-msa-d2f/web-script/s1.sh